# 7d2d/7d2dTimeReset/systemd/7d2d_timereset.timer
[Unit]
Description=Run 7DTD Time Reset every X minutes

[Timer]
OnCalendar=*-*-* *:*:00
Persistent=true


[Install]
WantedBy=timers.target
